# egpu-session-daemon@.service
[Unit]
Description=eGPU Session Selector for user %I
Before=display-manager.service
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/egpu-session-daemon %I
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
