# egpu-session-daemon@.service
[Unit]
Description=eGPU Session Selector for user %I
Before=display-manager.service
After=systemd-user-sessions.service

[Service]
Type=notify
ExecStart=/usr/local/bin/egpu-session-daemon %I
Restart=no
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
